# .dotfiles/.docker/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=__PROXY__"
